# strategies/scripts/com.tradepro.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job definition: refresh every ETF universe four times per
  day so the data behind the Compare page / digest is rarely more
  than a few hours stale. The 22:30 UTC slot is the canonical "after
  US close" run; the others fire around the open and close of the
  London + US sessions so a midday glance always sees data within
  the same trading day.

  Why 4× and not 1×: the audit (2026-05-05) noted the worker
  appearing offline because the heartbeat-staleness threshold is
  30 min, but the only refresh fires once per day. The user had to
  manually trigger compares to keep the cache fresh. Bumping the
  schedule to 4× covers a UK trader's natural check-in points
  without burning Yahoo quota — each run is ~60s on etf_uk_core,
  trivial vs the 6-hour gaps between fires.

  StartCalendarInterval with an array fires at every listed time:
    07:00 UTC — UK pre-open (08:00 BST)
    12:00 UTC — UK lunch / US pre-open
    17:00 UTC — UK close (18:00 BST), US mid-session
    22:30 UTC — after US close, definitive end-of-day

  Logs go to ~/.tradepro/logs/refresh-<date>.log.

  Install:

    bash strategies/scripts/install-launchd.sh

  The installer rewrites the {{HOME}} and {{PROJECT_DIR}} placeholders
  with the user's actual paths and copies the result to
  ~/Library/LaunchAgents/, then `launchctl load`s it. To pick up a
  schedule change, unload + reload:

    launchctl unload ~/Library/LaunchAgents/com.tradepro.refresh.plist
    launchctl load ~/Library/LaunchAgents/com.tradepro.refresh.plist

  Manual run (for testing):
    launchctl start com.tradepro.refresh

  Stop / unload:
    launchctl unload ~/Library/LaunchAgents/com.tradepro.refresh.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tradepro.refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>{{PROJECT_DIR}}/scripts/refresh.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>TRADEPRO_PROJECT_DIR</key>
    <string>{{PROJECT_DIR}}</string>
  </dict>

  <!-- Fire 4× / day. UTC times — adjust if you want them anchored to
       your local timezone instead. Skip times you don't care about by
       deleting the corresponding <dict>. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>7</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>17</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>22</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>{{HOME}}/.tradepro/logs/launchd-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>{{HOME}}/.tradepro/logs/launchd-stderr.log</string>

  <!-- Don't run if the laptop is closed or asleep at fire time;
       launchd will fire as soon as it wakes back up. -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
