# strategies/scripts/com.tradepro.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Persistent worker — alternative to com.tradepro.refresh when you
  want the Mac to stay continuously "alive" to the API instead of
  firing once-per-schedule and exiting.

  Differences vs com.tradepro.refresh:
    - KeepAlive: true       → launchd auto-restarts on crash / exit
    - RunAtLoad: true       → starts immediately on user login
    - Uses worker.sh which loops compare-runs every 30 min and
      heartbeats every 5 min in the background
    - No StartCalendarInterval (the inner loop owns the cadence)

  Pick ONE of refresh OR worker — running both schedules a redundant
  compare every 6 hours that the worker would have done anyway.
  install-launchd.sh defaults to the worker; add --refresh to install
  the cron-style schedule instead.

  Cadence overrides via the user's launchd env vars (set in this
  plist below if you want different defaults):
    WORKER_INTERVAL_SECONDS    default 1800 (30 min between cycles)
    HEARTBEAT_INTERVAL_SECONDS default 300  (5 min between pings)

  Pause without unloading:
    touch ~/.tradepro/worker.pause   # next cycle skips
    rm    ~/.tradepro/worker.pause   # resumes

  Stop completely:
    launchctl bootout gui/$UID ~/Library/LaunchAgents/com.tradepro.worker.plist

  Logs:
    ~/.tradepro/logs/worker-<date>.log
    ~/.tradepro/logs/worker-heartbeat-<date>.log
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tradepro.worker</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>{{PROJECT_DIR}}/scripts/worker.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>{{HOME}}</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>TRADEPRO_PROJECT_DIR</key>
    <string>{{PROJECT_DIR}}</string>
  </dict>

  <!-- Keep alive across crashes + restart on login. The inner loop
       handles its own pacing so we don't need launchd to care about
       throttling. -->
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>

  <!-- launchd's own throttle: don't restart faster than every 60s
       if the script crash-loops (defensive). -->
  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>{{HOME}}/.tradepro/logs/worker-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>{{HOME}}/.tradepro/logs/worker-stderr.log</string>
</dict>
</plist>
